# ansible-runner/roles/runner-manager/files/runner-manager.service.j2
[Unit]
Description=Runner manager (autoscaler) para runners efímeros de GitHub (org: TestingJDA)
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=5

Environment="GITHUB_ORG=TestingJDA"
Environment="GITHUB_TOKEN=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"
Environment="RUNNER_IMAGE=github-ephemeral-runner:latest"
Environment="MIN_RUNNERS=2"

ExecStart=/bin/bash -c '
  while true; do
    RUNNING=$(docker ps --filter "ancestor=${RUNNER_IMAGE}" -q | wc -l)
    RUNNING=$(echo $RUNNING | tr -d "[:space:]")
    if [ -z "${RUNNING}" ]; then RUNNING=0; fi

    if [ "${RUNNING}" -lt "${MIN_RUNNERS}" ]; then
      echo "[$(date)] Lanzando nuevo runner (actual ${RUNNING} < min ${MIN_RUNNERS})"
      docker run -d \
        --name runner_$(date +%s) \
        -e GITHUB_ORG=${GITHUB_ORG} \
        -e GITHUB_TOKEN=${GITHUB_TOKEN} \
        -v /var/run/docker.sock:/var/run/docker.sock \
        ${RUNNER_IMAGE}
    fi

    sleep 5
  done
'

ExecStop=/bin/bash -c '
  echo "Runner-manager detenido"
  docker ps --filter "ancestor=${RUNNER_IMAGE}" -q | xargs -r docker stop
'

[Install]
WantedBy=multi-user.target